# internal/render/testdata/golden/systemd/user/bossa.target
# Generated by stompbox — DO NOT EDIT.
[Unit]
Description=Stompbox session "bossa"
Wants=carla@clean.service hydrogen@bossa_groove.service sooperlooper@looper.service qpwgraph@bossa.service
