# internal/render/testdata/golden/systemd/user/carla@.service
# Generated by stompbox — DO NOT EDIT.
[Unit]
Description=Carla instance "%i" (stompbox)
# Host headless : pas de BindsTo=graphical-session.target,
# l'IHM (carla-control) est un client jetable lancé via "stomp gui".
After=pipewire.service
# Requires= ne propage pas l'arrêt : c'est ce flag qui fait que
# "stomp off" (stop du target) arrête réellement l'instance.
StopWhenUnneeded=yes
# Restart=always sans limite : seul "stomp off" éteint le pédalier.
StartLimitIntervalSec=0

[Service]
Type=simple
EnvironmentFile=%h/.config/stompbox/carla/%i.env
# pw-jack : la libjack système est celle de JACK2 (sans jackd) ;
# il redirige vers la libjack de PipeWire.
ExecStart=/usr/bin/pw-jack /usr/bin/carla --no-gui ${CARLA_PRESET}
Restart=always
RestartSec=2
