# internal/render/testdata/golden/systemd/user/hydrogen@.service
# Generated by stompbox — DO NOT EDIT.
[Unit]
Description=Hydrogen instance "%i" (stompbox)
After=pipewire.service graphical-session.target
BindsTo=graphical-session.target
StopWhenUnneeded=yes

[Service]
Type=simple
EnvironmentFile=%h/.config/stompbox/hydrogen/%i.env
# pw-jack : la libjack système est celle de JACK2 (sans jackd) ;
# il redirige vers la libjack de PipeWire.
ExecStart=/usr/bin/pw-jack /usr/bin/hydrogen -n --driver jack -s ${H2_SONG}
Restart=on-failure
RestartSec=2
