# internal/render/testdata/golden/systemd/user/jam.target
# Generated by stompbox — DO NOT EDIT.
[Unit]
Description=Stompbox session "jam"
Wants=carla@dirty.service sooperlooper@jam_looper.service qpwgraph@jam.service
