# internal/render/testdata/golden/systemd/user/qpwgraph@.service
# Generated by stompbox — DO NOT EDIT.
[Unit]
Description=qpwgraph patchbay "%i" (stompbox)
After=pipewire.service graphical-session.target
BindsTo=graphical-session.target
StopWhenUnneeded=yes

[Service]
Type=simple
EnvironmentFile=%h/.config/stompbox/qpwgraph/%i.env
# -a : patchbay activé, les connexions sont réappliquées dès que les
# nœuds apparaissent — l'ordre de démarrage des apps devient indifférent.
ExecStart=/usr/bin/qpwgraph -m -a $QPW_EXCLUSIVE ${QPW_PATCHBAY}
Restart=on-failure
RestartSec=2
