# internal/render/testdata/golden/systemd/user/sooperlooper@.service
# Generated by stompbox — DO NOT EDIT.
[Unit]
Description=SooperLooper engine instance "%i" (stompbox)
# Moteur headless : pas de BindsTo=graphical-session.target,
# l'IHM (slgui) est un client jetable lancé via "stomp gui".
After=pipewire.service
StopWhenUnneeded=yes

[Service]
Type=simple
EnvironmentFile=%h/.config/stompbox/sooperlooper/%i.env
# -j sl_%i : nom de client JACK déterministe, requis pour que le
# patchbay qpwgraph s'applique de façon fiable.
# pw-jack : la libjack système est celle de JACK2 (sans jackd) ;
# il redirige vers la libjack de PipeWire.
ExecStart=/usr/bin/pw-jack /usr/bin/sooperlooper -q -j sl_%i -p $SL_OSC_PORT -l $SL_LOOPS -c $SL_CHANNELS -t $SL_LOOPTIME $SL_EXTRA_ARGS
Restart=on-failure
RestartSec=2
